# conf/crontab.cron
#*/30, * * * * /home/pi/code/jfkspoken/tweet.sh
0,15,30,45 * * * * /Users/pr/work/code/jfkspoken/tweet.sh
